# mint/90-offline-charging.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", ACTION=="change", TEST=="/tmp/charger_mode", RUN+="/usr/bin/poweroff -f"
